# MF761 Series/MF761C/files/recovery/etc/inittab
::sysinit:/recovery/etc/rc
::respawn:-/recovery/bin/login

