# R186X Series/R1869-WIFI4/files/etc/inittab
::sysinit:/etc/rc
::respawn:-/bin/login

